# jeffsender-queue.conf
[program:jeffsender-queue]
; ============================================================
; Supervisor Config untuk Jeffsender Queue Worker
; Simpan file ini di: /etc/supervisor/conf.d/jeffsender-queue.conf
;
; Kemudian jalankan:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start jeffsender-queue
; ============================================================

command=/usr/bin/php /var/www/jeffsender/artisan queue:work database --timeout=900 --tries=1 --sleep=3 --max-jobs=50
directory=/var/www/jeffsender
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/jeffsender/storage/logs/queue-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=600
